# regranbill-backend.service
[Unit]
Description=ReGranBill Backend
After=network.target postgresql.service

[Service]
Type=simple
WorkingDirectory=/var/www/regranbill/backend-publish
ExecStart=/usr/bin/dotnet /var/www/regranbill/backend-publish/ReGranBill.Server.dll
Restart=always
RestartSec=10
User=www-data
Group=www-data

# Environment
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://0.0.0.0:5105
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=regranbill-backend

[Install]
WantedBy=multi-user.target
